# ddclient/rootfs/etc/s6-overlay/s6-rc.d/ha-ddclient/run
#!/command/with-contenv bashio
# ==============================================================================
# Start the ddclient service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

if [ -f /share/ddclient/ddclient.conf ]; then
    bashio::log.error "Missing /share/ddclient/ddclient.conf. Add file and restart addon."
    exit
fi

## Print the message the user supplied
bashio::log.info "Starting ddclient addon"

## Run your program
exec /usr/bin/ddclient -verbose -foreground -file /share/ddclient/ddclient.conf
