# ddclient/rootfs/etc/s6-overlay/s6-rc.d/ddclient/run
#!/command/execlineb -P

/usr/bin/ddclient -foreground -file /share/ddclient/ddclient.conf